# crontab
CRON_TZ=Europe/Berlin

*/1 * * * * xh http://host.docker.internal:8081
